# utils/crontab/crontab
0 3 * * * /home/drag0n/.bin/get_exchange_rate.sh
*/30 * * * * /home/drag0n/.i3/scripts/get_weather.sh
