# utils/crontab/crontab
*/10 * * * * /home/drag0n/.bin/get_mail.sh

