# utils/crontab/crontab
*/10 * * * * /usr/bin/env fdm -l fetch

